# config/swictation-daemon.service
[Unit]
Description=Swictation Voice-to-Text Daemon
After=graphical-session.target
# Auto-start UI service when daemon starts
Wants=swictation-ui.service

[Service]
Type=simple
ExecStart=/usr/local/lib/node_modules/swictation/bin/swictation-daemon
Restart=on-failure
RestartSec=5

# Environment - CRITICAL for GPU acceleration and correct model operation
Environment="RUST_LOG=info"

# ORT_DYLIB_PATH: Points to ONNX Runtime library (CRITICAL - without this, model produces blank output)
Environment="ORT_DYLIB_PATH=/usr/local/lib/node_modules/swictation/lib/native/libonnxruntime.so"

# LD_LIBRARY_PATH: CUDA runtime + ONNX Runtime CUDA providers (for GPU acceleration)
# Without this: 14x performance regression (43s vs 3-6s expected for 57s audio)
Environment="LD_LIBRARY_PATH=/usr/local/cuda/lib64:/usr/local/lib/node_modules/swictation/lib/native"

# CUDA_HOME: Helps CUDA libraries find additional resources at runtime
Environment="CUDA_HOME=/usr/local/cuda"

# Import user environment for PulseAudio/PipeWire session (enables full audio device detection)
# This ensures all audio devices are detected properly (4 devices instead of 1)
ImportEnvironment=

# Resource limits - REMOVED for development/testing
# Previously: MemoryMax=500M caused 99.7% usage + 2.4GB swap
# Previously: CPUQuota=50% limited to 4 cores
# 1.1B model needs ~3.5GB GPU + ~1GB system RAM + full CPU for optimal performance
# Uncomment and adjust for production if needed:
# MemoryMax=8G
# CPUQuota=200%

# Security hardening (relaxed for USB device access)
#PrivateTmp=true
#ProtectSystem=strict
#ProtectHome=read-only
#ReadWritePaths=/tmp
#ReadWritePaths=%h/.local/share/swictation
#ReadWritePaths=%h/.cache/swictation
#ReadWritePaths=%h/.config/swictation

[Install]
WantedBy=default.target
