# config/swictation.service
[Unit]
Description=Swictation voice dictation daemon
Documentation=https://github.com/yourusername/swictation
After=graphical-session.target
PartOf=graphical-session.target
Wants=swictation-tray.service

[Service]
Type=simple
# Use Python 3.12 if available, otherwise system Python 3
ExecStart=/bin/bash -c 'if command -v python3.12 &> /dev/null; then python3.12 /opt/swictation/src/swictationd.py; else /usr/bin/env python3 /opt/swictation/src/swictationd.py; fi'
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment
Environment="PYTHONUNBUFFERED=1"
Environment="PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin"

# Resource limits
MemoryMax=6G
CPUQuota=200%

# Security hardening
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
NoNewPrivileges=true

[Install]
WantedBy=sway-session.target
# Fallback to default.target if sway-session.target doesn't exist
Also=default.target
# Ensure tray service is also enabled/started
Also=swictation-tray.service
