# config/swictation-ui.service
[Unit]
Description=Swictation System Tray UI
# Start after daemon is ready
After=swictation-daemon.service
# Required dependency - UI needs daemon to function
Requires=swictation-daemon.service
# Stop when daemon stops
BindsTo=swictation-daemon.service
# Part of daemon service group
PartOf=swictation-daemon.service

[Service]
Type=simple
ExecStart=/usr/local/lib/node_modules/swictation/bin/swictation-ui-launcher
Restart=on-failure
RestartSec=5

# Environment for Wayland/X11
Environment="WAYLAND_DISPLAY=wayland-1"
Environment="XDG_RUNTIME_DIR=/run/user/1000"
Environment="QT_QPA_PLATFORM=wayland"

# Resource limits (UI is lightweight)
MemoryMax=100M
CPUQuota=10%

[Install]
WantedBy=default.target
