# config/swictation-tauri.service
[Unit]
Description=Swictation Tauri UI
After=swictation-daemon.service graphical-session.target
Requires=swictation-daemon.service

[Service]
Type=simple
ExecStart=/opt/swictation/tauri-ui/src-tauri/target/release/swictation-ui
Restart=on-failure
RestartSec=5
Environment="DISPLAY=:0"
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
