# config/swictation-daemon.service
[Unit]
Description=Swictation Voice-to-Text Daemon
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/local/lib/node_modules/swictation/bin/swictation-daemon
Restart=on-failure
RestartSec=5

# Environment
Environment="RUST_LOG=info"

# Resource limits
MemoryMax=500M
CPUQuota=50%

# Security hardening
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/tmp
ReadWritePaths=%h/.local/share/swictation
ReadWritePaths=%h/.cache/swictation
ReadWritePaths=%h/.config/swictation

[Install]
WantedBy=default.target