# config/swictation.service
[Unit]
Description=Swictation voice dictation daemon (Rust)
Documentation=https://github.com/ruvnet/swictation
After=graphical-session.target
PartOf=graphical-session.target
Wants=swictation-tauri.service

[Service]
Type=simple
# Rust daemon binary
ExecStart=/opt/swictation/bin/swictationd
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment
Environment="RUST_LOG=info"
Environment="PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin"

# Resource limits (GPU transcription can be memory-intensive)
MemoryMax=8G
CPUQuota=200%

# Security hardening
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
NoNewPrivileges=true
# Allow access to GPU devices
DeviceAllow=/dev/nvidia* rw
DeviceAllow=/dev/dri/* rw

[Install]
WantedBy=default.target
# Ensure Tauri UI service is also enabled/started
Also=swictation-tauri.service
